# deploy/systemd/skoljka.service.template
[Unit]
Description=Skoljka Django application
After=network.target postgresql.service

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__REPO_ROOT__
EnvironmentFile=-/etc/skoljka/skoljka.env
ExecStart=__REPO_ROOT__/venv/bin/gunicorn skoljka.config.wsgi:application --bind 127.0.0.1:__APP_PORT__ --workers 3 --timeout 60 --access-logfile - --error-logfile -
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
